# tcs-connect/tcs-connect.service
[Unit]
Description=TCS Connect
Requires=sshd.service

[Service]
Type=simple
ExecStart=/home/tsm/tcs-connect/start.sh
Restart=always
RestartSec=3
